# deploy/synapse.service
[Unit]
Description=Synapse LLM Controller
Documentation=https://github.com/MacCracken/synapse
After=network-online.target
Wants=network-online.target

[Service]
Type=notify
ExecStart=/usr/local/bin/synapse serve
Environment=SYNAPSE_CONFIG=/etc/synapse/synapse.toml
Environment=RUST_LOG=info

# Run as dedicated service user
User=synapse
Group=synapse

# Restart policy
Restart=on-failure
RestartSec=5
WatchdogSec=30
ExecHealthCheck=/usr/bin/curl -sf http://127.0.0.1:8420/health

# Resource limits
LimitNOFILE=65536
LimitMEMLOCK=infinity

# Security hardening
# NOTE: Drop-in overrides for inference-only and training modes are available.
# See deploy/synapse-inference.conf and deploy/synapse-training.conf.
# Install to /etc/systemd/system/synapse.service.d/ as needed.
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
NoNewPrivileges=true
ReadWritePaths=/var/lib/synapse
ReadOnlyPaths=/etc/synapse

# GPU access — NVIDIA devices and DRI for AMD ROCm
DeviceAllow=/dev/nvidia* rwm
DeviceAllow=/dev/nvidiactl rwm
DeviceAllow=/dev/nvidia-uvm rwm
DeviceAllow=/dev/nvidia-uvm-tools rwm
DeviceAllow=/dev/dri/* rwm

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=synapse

[Install]
WantedBy=multi-user.target
